# lib/systemd/system/foo.service
[Unit]
Description=Foo

[Service]
Type=simple
ExecStart=/usr/local/bin/foo.executable

[Install]
WantedBy=multi-user.target
